# games/owasp_juice_shop/supervisord.conf
[program:owasp_juice_shop]
command= npm start
directory=/juice-shop
autostart=true
autorestart=true
