# games/owasp_juice_shop/supervisord.conf
[program:owasp_juice_shop]
command=node build/app.js
directory=/juice-shop
autostart=true
user=juice-shop
autorestart=true
